# cf_rest_updater.service
[Unit]
Description=Cloudflare DNS updater service

[Service]
Type=oneshot
ExecStart=/your/dir/here/updater.sh
